# imageroot/systemd/user/workspace@.service
[Unit]
Description=Hermes Workspace sidecar %i
Requires=hermes-pod@%i.service hermes@%i.service
After=hermes-pod@%i.service hermes@%i.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Environment=HERMES_AGENT_WORKSPACE_IMAGE=ghcr.io/nethserver/hermes-agent-workspace:latest
EnvironmentFile=-%S/state/environment
EnvironmentFile=-%S/state/agent_%i.env
EnvironmentFile=-%S/state/agent_%i_secrets.env
WorkingDirectory=%S/state
Restart=on-failure
KillMode=none
Type=notify
NotifyAccess=all
ExecStartPre=/usr/local/bin/runagent discover-smarthost
ExecStartPre=/usr/local/bin/runagent sync-agent-runtime --agent-id %i
ExecStartPre=-/usr/bin/podman rm --force workspace-%i
ExecStart=runagent /usr/bin/podman run \
		--name workspace-%i \
		--replace \
		--rm \
		--sdnotify=conmon \
		--cgroups=no-conmon \
		--pod hermes-pod-%i \
		--env-file %S/state/agent_%i.env \
		--env-file %S/state/agent_%i_secrets.env \
		--env HOST=127.0.0.1 \
		--env PORT=3000 \
		--env TRUST_PROXY=1 \
		--env HERMES_API_URL=http://127.0.0.1:8642 \
		--env HERMES_API_TOKEN=${API_SERVER_KEY} \
		--env HERMES_DASHBOARD_URL=http://127.0.0.1:9120 \
		--env HOME=/opt/data \
		--env HERMES_HOME=/opt/data/.hermes \
		--env HERMES_WORKSPACE_DIR=/opt/data/workspace \
		--volume hermes-agent-%i-home:/opt/data:z \
		--tz=${TIMEZONE} \
		${HERMES_AGENT_WORKSPACE_IMAGE}
ExecStop=/usr/bin/podman stop --ignore --time 10 workspace-%i
ExecStopPost=-/usr/bin/podman rm --force workspace-%i
TimeoutStopSec=70

[Install]
WantedBy=default.target